# roles/rootio/files/supervisor/rootio_web.conf
[program:rootio_web]
user = rootio
command = /opt/rootio/bin/rootio_web runserver
redirect_stderr = true
