# roles/rootio/files/supervisor/station_telephony.conf
[program:station_telephony]
user = rootio
command = /opt/rootio/bin/station-telephony-run
redirect_stderr = true
